# config/launchd/com.fohte.headroom-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.fohte.headroom-proxy</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/fohte/.local/bin/mise</string>
    <string>x</string>
    <string>--</string>
    <string>headroom</string>
    <string>proxy</string>
    <string>--port</string>
    <string>8787</string>
    <string>--mode</string>
    <string>cache</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/fohte/.local/state/headroom/proxy.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/fohte/.local/state/headroom/proxy.log</string>
</dict>
</plist>
