# config/launchd/com.fohte.skk-watcher.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.fohte.skk-watcher</string>

  <key>WatchPaths</key>
  <array>
    <string>__HOME__/Library/Application Support/AquaSKK/skk-jisyo.utf8</string>
  </array>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>__SKK_DICT_REPO__/scripts/skk-watch</string>
    <string>--local</string>
    <string>__HOME__/Library/Application Support/AquaSKK/skk-jisyo.utf8</string>
    <string>--snapshot</string>
    <string>__SKK_DICT_REPO__/.git/skk-sync-snapshot</string>
    <string>--events</string>
    <string>__SKK_DICT_REPO__/.git/skk-events.jsonl</string>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>__HOME__/.local/state/skk-watcher/watch.log</string>

  <key>StandardErrorPath</key>
  <string>__HOME__/.local/state/skk-watcher/watch.log</string>
</dict>
</plist>
